# set-hostname.service
[Unit]
Description=IP address as hostname
Requires=network-online.target
After=network-online.target vpn-client.service

[Service]
ExecStart=/bin/bash /bin/set-hostname.sh
RemainAfterExit=yes

[Install]
WantedBy=basic.target